# install/files/inotifywatcher.service
[Unit]
Description=inotifywatcher

[Service]
ExecStart=/data/admin/scripts/inotifywatcher.sh

[Install]
WantedBy=multi-user.target
